# etc/service/dash/run
#!/bin/bash

export LC_ALL=C.UTF-8
export LANG=C.UTF-8

exec /dash.py 2>&1 | logger -t $DASH_LOG_FACILITY
